# warden_archiver.cron
# gzip old data in the warden archive every night (at 00:10)
# TODO: maybe call directly by archiver?
10 0 * * * wardenarchiver /data/warden_archive/process_daily_file.sh $(date -Idate -d "-1day") >>/data/warden_archive/processing.log 2>&1
